# example/farside_bot.timer
[Unit]
Description="Timer to run The Far Side Comic Telegram Bot daily"

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=multi-user.target
